# etc/systemd/system/docker.service.d/override.conf
# you can create this file and the folder structure by executing "sudo systemctl edit docker"
[Service]
ExecStartPre=/usr/local/bin/wait-for-unlock.sh
