# etc/systemd/system/docker.service.d/override.conf
[Service]
ExecStartPre=/usr/local/bin/wait-for-unlock.sh
